# luci-theme-uniwrt/root/etc/uci-defaults/30_luci-theme-uniwrt
#!/bin/sh
# luci-theme-uniwrt — register the theme and make it the active LuCI theme.
# Runs once on install; uci-defaults removes the script afterwards.

uci -q batch <<-EOF
	set luci.themes.UniWRT=/luci-static/uniwrt
	set luci.main.mediaurlbase=/luci-static/uniwrt
	commit luci
EOF

# Clear any cached compiled templates so the new theme is picked up immediately.
rm -f /tmp/luci-indexcache /tmp/luci-modulecache/* 2>/dev/null
rm -rf /tmp/luci-modulecache 2>/dev/null

exit 0
